# crontab
*/5 * * * * /Users/olamideabioro/Documents/data analytics project/city_mobility_ops/scripts/cron_run_pipeline.sh
